# supervisord.conf
#Supervisord is a simple and popular choice for process monitoring. We are using this so that a single line of supervisord should take care of starting all the related services.


[supervisord]
nodaemon=true

[program:apache]
command=/bin/bash -c "exec /usr/sbin/apachectl -DFOREGROUND"

[program:nagios]
command=/bin/bash -c "/usr/local/nagios/bin/nagios /usr/local/nagios/etc/nagios.cfg"

[program:nagios-bugfix]
command=/bin/bash -c "sleep 10; chmod 666 /usr/local/nagios/var/rw/nagios.cmd; cat;"

[program:sendmail]
command=/bin/bash -c "exec /etc/init.d/sendmail start"
